# autostart/mount_RAW_DATA
#!/bin/bash
### BEGIN INIT INFO
# Provides: RAW_DATA
# Required-Start: $local_fs
# Required-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Mount RAW_DATA partition to /media/RAW_DATA
### END INIT INFO
mount /dev/mmcblk0p1 /media/RAW_DATA -o dmask=000,fmask=111
exit
